# .config/containers/systemd/soularr.container
[Unit]
Description=Soularr Lidarr-to-Soulseek bridge
Wants=network-online.target lidarr.service slskd.service
After=network-online.target lidarr.service slskd.service

[Container]
ContainerName=soularr
Image=docker.io/mrusse08/soularr:latest
Pull=newer
Network=media-stack.network
UserNS=keep-id
PodmanArgs=--user=1000:1000 --cpu-shares=128
PublishPort=127.0.0.1:8265:8265
Environment=PUID=1000
Environment=PGID=1000
Environment=TZ=America/New_York
Environment=SCRIPT_INTERVAL=300
Volume=/var/home/evan/.local/share/media-stack/soularr:/data:Z
Volume=/mnt/Media/Downloads:/downloads:z

[Service]
Restart=on-failure
RestartSec=30
Nice=10
CPUWeight=20
IOWeight=20
IOSchedulingClass=idle

[Install]
WantedBy=default.target
